# scripts/eve-pi-manager-ws.service
[Unit]
Description=EVE PI Manager - zKillboard R2Z2 Poller
After=rabbitmq-server.service postgresql.service network.target
Wants=rabbitmq-server.service postgresql.service

[Service]
Type=simple
User=evepi
Group=evepi
WorkingDirectory=/opt/eve-pi-manager
EnvironmentFile=/opt/eve-pi-manager/.env
Environment=CELERY_WS_AUTOSTART=1
ExecStart=/opt/eve-pi-manager/venv/bin/celery -A app.celery_app worker --loglevel=info --concurrency=1 -Q ws --hostname=ws@%H
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eve-pi-manager-ws
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/eve-pi-manager

[Install]
WantedBy=multi-user.target
