# scripts/planetflow-ws.service
[Unit]
Description=PlanetFlow - zKillboard R2Z2 Poller
After=rabbitmq-server.service postgresql.service network.target
Wants=rabbitmq-server.service postgresql.service

[Service]
Type=simple
User=planetflow
Group=planetflow
WorkingDirectory=/opt/planetflow
EnvironmentFile=/opt/planetflow/.env
ExecStart=/opt/planetflow/venv/bin/celery -A app.celery_app worker --loglevel=info --concurrency=1 -Q ws --hostname=ws@%H
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=planetflow-ws
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/planetflow

[Install]
WantedBy=multi-user.target
