# scripts/planetflow.service
[Unit]
Description=PlanetFlow - Web
After=postgresql.service rabbitmq-server.service network.target
Wants=postgresql.service rabbitmq-server.service

[Service]
Type=simple
User=planetflow
Group=planetflow
WorkingDirectory=/opt/planetflow
EnvironmentFile=/opt/planetflow/.env
ExecStart=/opt/planetflow/venv/bin/gunicorn app.main:app \
    -k uvicorn.workers.UvicornWorker \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=planetflow
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/planetflow

[Install]
WantedBy=multi-user.target
